# catfit.service
[Unit]
Description=Cat weight measuring
After=multi-user.target

[Service]
User=pi
Group=pi
WorkingDirectory=/home/pi/git/catfit
ExecStart=/home/pi/git/catfit/venv/bin/python /home/pi/git/catfit/catfit.py

[Install]
WantedBy=multi-user.target
